# deploy/ally-watchdog.service
[Unit]
Description=Ally Watchdog (Timeout Scanner)
After=network.target docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ally
ExecStart=/opt/ally/bin/ally watch
Restart=on-failure
RestartSec=5
EnvironmentFile=/opt/ally/.env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ally-watchdog

[Install]
WantedBy=multi-user.target
